# novnc-display/services/xfce-session/run
#!/bin/sh
export DISPLAY=:0
export XDG_RUNTIME_DIR=/tmp/runtime-guest
# Wait for Xvfb to be ready
while ! xset -display :0 q > /dev/null 2>&1; do sleep 0.1; done
# Drop privileges and start session
exec chpst -u guest dbus-run-session -- bash -c 'xfconf-query -c xsettings -p /Net/IconThemeName -n -t string -s "Humanity" || true; startxfce4'
