# novnc-display/services/x11vnc/run
#!/bin/sh
exec x11vnc -display :0 -nopw -forever -xkb -q
